# deploy/cron/weekly_walkforward_governance.cron
0 6 * * 1 cd /workspaces/btc && PYTHONPATH=/workspaces/btc /workspaces/btc/.venv/bin/python -m src.scripts.run_walkforward_weekly >> logs/run_walkforward_weekly.log 2>&1
